# src/leggedrobotics-anymal_research_rsl-70c2caecb930/anymal/anymal_b/anymal_b_drivers/anymal_b_pc/install/etc/profile.d/anymal_catkin_ws_overlay.sh
#!/bin/bash

source /home/integration/catkin_ws/devel/setup.sh
