# src/leggedrobotics-anymal_research_rsl-70c2caecb930/anymal/anymal_b/anymal_b_drivers/anymal_b_pc/install/etc/cron.d/check_disk_space
# Check disk space every 10 minutes
*/10 * * * * integration /home/integration/git/anymal_b_drivers/anymal_b_pc/bin/check_disk_space.sh /home/integration/.ros
